# home/dot_config/systemd/user/resticprofile@home.timer
[Unit]
Description=Hourly resticprofile %i backup timer

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
